# systemd/docker-at.service.in
[Unit]
Description=%I docker container
Requires=docker.service
After=network.target docker.service

[Service]
Type=oneshot
ExecStart=@LIBEXECDIR@/systemd-docker-wrapper.sh -c %i start
ExecStop=@LIBEXECDIR@/systemd-docker-wrapper.sh -c %i stop
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
